# docker/sbdoc-with-mongo-alpine/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod

[program:doclever]
command=/usr/bin/node /root/DOClever/Server/bin/www
