# docker/sbdoc-with-mongo/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod

[program:DOClever]
command=/usr/bin/node /root/DOClever/Server/bin/www
